# debian/castor-vmgr-server.postinst
/sbin/chkconfig --add vmgrdaemon
/sbin/service vmgrdaemon condrestart > /dev/null 2>&1 || :
