# debian/castor-dlf-server.postinst
/sbin/chkconfig --add dlfserver
/sbin/service dlfserver condrestart > /dev/null 2>&1 || :
