# debian/castor-vdqm2-server.postinst
/sbin/chkconfig --add vdqmserver
/sbin/service vdqmserver condrestart > /dev/null 2>&1 || :
