# debian/castor-tape-server.postinst
/sbin/chkconfig --add tpdaemon
/sbin/service tpdaemon condrestart > /dev/null 2>&1 || :
