# debian/castor-repack-server.postinst
/sbin/chkconfig --add repackserver
/sbin/service repackserver condrestart > /dev/null 2>&1 || :
