# debian/castor-ns-server.postinst
/sbin/chkconfig --add nsdaemon
/sbin/service nsdaemon condrestart > /dev/null 2>&1 || :
