# debian/castor-rmc-server.postinst
/sbin/chkconfig --add rmcdaemon
/sbin/service rmcdaemon condrestart > /dev/null 2>&1 || :
