# debian/castor-repack-server.postinst
/sbin/chkconfig --add repackd
/sbin/service repackd condrestart > /dev/null 2>&1 || :
