# debian/castor-rmnode-server.postinst
/sbin/chkconfig --add rmNodeDaemon
/sbin/service rmNodeDaemon condrestart > /dev/null 2>&1 || :
