# debian/castor-tape-server.postinst
/sbin/chkconfig --add taped
/sbin/service taped condrestart > /dev/null 2>&1 || :
