# debian/castor-mighunter-server.postinst
/sbin/chkconfig --add mighunter
/sbin/service mighunter condrestart > /dev/null 2>&1 || :
