# debian/castor-gc-server.postinst
/sbin/chkconfig --add gcDaemon
/sbin/service gcDaemon condrestart > /dev/null 2>&1 || :
/bin/chown stage:st /var/spool/gc
/bin/chown stage:st /var/spool/gc/log > /dev/null 2>&1
