# debian/castor-repack-server.postinst
/sbin/chkconfig --add repackd
/sbin/service repackd condrestart > /dev/null 2>&1 || :
/bin/chown stage:st /var/spool/repack
/bin/chown stage:st /var/spool/repack/log > /dev/null 2>&1
