# debian/castor-tapebridge-server.postinst
/sbin/chkconfig --add tapebridged
/sbin/service tapebridged condrestart > /dev/null 2>&1 || :
