# debian/castor-gc-server.postinst
/sbin/chkconfig --add gcd
/sbin/service gcd condrestart > /dev/null 2>&1 || :
/bin/chown stage:st /var/spool/gc
/bin/chown stage:st /var/spool/gc/log > /dev/null 2>&1
