# debian/castor-rtcopy-server.postinst
/sbin/chkconfig --add rtcpd
/sbin/service rtcpd condrestart > /dev/null 2>&1 || :
