# debian/castor-stager-server.postinst
/sbin/chkconfig --add stager
/sbin/service stager condrestart > /dev/null 2>&1 || :
