# debian/castor-ns-server.postinst
/sbin/chkconfig --add nsd
/sbin/service nsd condrestart > /dev/null 2>&1 || :
