# debian/castor-rtcopy-clientserver.postinst
/sbin/chkconfig --add rtcpclientd
/sbin/service rtcpclientd condrestart > /dev/null 2>&1 || :
