# debian/castor-rechandler-server.postinst
/sbin/chkconfig --add rechandler
/sbin/service rechandler condrestart > /dev/null 2>&1 || :
