# debian/castor-rmnode-server.postinst
/sbin/chkconfig --add rmnoded
/sbin/service rmnoded condrestart > /dev/null 2>&1 || :
