# debian/castor-tapegateway-server.postinst
/sbin/chkconfig --add tapegatewayd
/sbin/service tapegatewayd condrestart > /dev/null 2>&1 || :
