# debian/castor-tapeserver-server.postinst
/sbin/chkconfig --add tapeserverd
/sbin/service tapeserverd condrestart > /dev/null 2>&1 || :
