# debian/castor-rechandler-server.postinst
/sbin/chkconfig --add rechandlerd
/sbin/service rechandlerd condrestart > /dev/null 2>&1 || :
