# debian/castor-vdqm-server.postinst
/sbin/chkconfig --add vdqmserv
/sbin/service vdqmserv condrestart > /dev/null 2>&1 || :
