# debian/castor-ns-server.postinst
/usr/sbin/groupadd -g 1474 st > /dev/null 2>&1 || :
/usr/sbin/useradd -s /bin/nologin -c "Stager Staging System" -u 14029 -g 1474 stage > /dev/null 2>&1 || :
/sbin/chkconfig --add nsd
/sbin/service nsd condrestart > /dev/null 2>&1 || :
