# debian/castor-upv-server.postinst
/sbin/chkconfig --add Cupvdaemon
/sbin/service Cupvdaemon condrestart > /dev/null 2>&1 || :
