# debian/castor-vmgr-server.postinst
/sbin/chkconfig --add vmgrd
/sbin/service vmgrd condrestart > /dev/null 2>&1 || :
