# debian/castor-cleaning-server.postinst
/sbin/chkconfig --add cleaning
/sbin/service cleaning condrestart > /dev/null 2>&1 || :
