# debian/castor-gc-server.postinst
/sbin/chkconfig --add gcDaemon
/sbin/service gcDaemon condrestart > /dev/null 2>&1 || :
