# debian/castor-vdqm2-server.postinst
/sbin/chkconfig --add vdqmd
/sbin/service vdqmd condrestart > /dev/null 2>&1 || :
