# debian/castor-rh-server.postinst
/sbin/chkconfig --add rhserver
/sbin/service rhserver condrestart > /dev/null 2>&1 || :
