# debian/castor-aggregator-server.postinst
/sbin/chkconfig --add aggregatord
/sbin/service aggregatord condrestart > /dev/null 2>&1 || :
