# CTA/tapeserver/castor/debian/castor-transfer-manager.postinst
/sbin/chkconfig --add transfermanagerd
/sbin/service transfermanagerd condrestart > /dev/null 2>&1 || :
