# CTA/tapeserver/castor/debian/castor-traffic-shaper.postinst
/sbin/chkconfig --add castor-traffic-shaper
/sbin/service castor-traffic-shaper start > /dev/null 2>&1 || :
